# environment/release/systemd/illuvi-analytics-backend-watcher.path
[Path]
Unit=illuvi-analytics-backend-watcher.service
PathModified=/opt/backend

[Install]
WantedBy=multi-user.target